# tmp/com.agentcore.beadmonitor.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.agentcore.beadmonitor</string>
	<key>ProgramArguments</key>
	<array>
		<string>/bin/bash</string>
		<string>/Users/james/Projects/AgentCore/flywheel_tools/scripts/beads/bead-stale-monitor.sh</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>KeepAlive</key>
	<true/>
	<key>ThrottleInterval</key>
	<integer>10</integer>
	<key>StandardOutPath</key>
	<string>/Users/james/Projects/AgentCore/tmp/launchd-beadmonitor.log</string>
	<key>StandardErrorPath</key>
	<string>/Users/james/Projects/AgentCore/tmp/launchd-beadmonitor.err</string>
	<key>WorkingDirectory</key>
	<string>/Users/james/Projects/AgentCore</string>
	<key>EnvironmentVariables</key>
	<dict>
		<key>PATH</key>
		<string>/Users/james/.local/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
	</dict>
</dict>
</plist>
